# lib/wifi_page.service
[Unit]
Description=start flask server for wifi config
 
[Service]
Type=simple
ExecStart=/usr/bin/python /home/pi/iotchip/bin/server/wifi.py
Restart=on-abort
 
[Install]
WantedBy=multi-user.target
